# dhan-algo.service
[Unit]
Description=Dhan Options Auto-Trader
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/dhan-algo-trade
ExecStart=/home/ubuntu/dhan-algo-trade/venv/bin/python /home/ubuntu/dhan-algo-trade/dhan_live.py

# Restart settings (important for trading bot)
Restart=always
RestartSec=5

# Logs
StandardOutput=append:/home/ubuntu/dhan-algo-trade/dhan.log
StandardError=append:/home/ubuntu/dhan-algo-trade/dhan-error.log

# Env
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target